# roles/factorio/templates/update-factorio.timer.j2
[Unit]
Description=Update factorio server

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
